# install_files/udev/99-razer-wolverine.rules
# Razer Wolverine V3 Pro - Automatically switch from xpad to razerwolverine
# These rules rebind the devices to razerwolverine after they are detected

# Wolverine V3 Pro for Xbox (USB-C) - Interface 0
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1532", ATTRS{idProduct}=="0a3f", ATTR{bInterfaceNumber}=="00", \
    RUN+="/usr/lib/udev/razer-wolverine-rebind.sh %k"

# Wolverine V3 Pro 8K PC Wired (USB-C) - Interface 0
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1532", ATTRS{idProduct}=="0a57", ATTR{bInterfaceNumber}=="00", \
    RUN+="/usr/lib/udev/razer-wolverine-rebind.sh %k"

# Wolverine V3 Pro 8K PC Wireless (2.4GHz Dongle) - Interface 0
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1532", ATTRS{idProduct}=="0a59", ATTR{bInterfaceNumber}=="00", \
    RUN+="/usr/lib/udev/razer-wolverine-rebind.sh %k"
